# Semester_A/01 Extract data from MITRE ATTACK/Persistence/Persistence_part_5.json
{
    "T1556.008": [
        "Network Provider DLL",
        "Adversaries may register malicious network provider dynamic link libraries (DLLs) to capture cleartext user credentials during the authentication process"
    ],
    "T1137": [
        "Office Application Startup",
        "Adversaries may leverage Microsoft Office-based applications for persistence between startups"
    ],
    "T1137.001": [
        "Office Template Macros",
        "Adversaries may abuse Microsoft Office templates to obtain persistence on a compromised system"
    ],
    "T1137.002": [
        "Office Test",
        "Adversaries may abuse the Microsoft Office \"Office Test\" Registry key to obtain persistence on a compromised system"
    ],
    "T1137.003": [
        "Outlook Forms",
        "Adversaries may abuse Microsoft Outlook forms to obtain persistence on a compromised system"
    ],
    "T1137.004": [
        "Outlook Home Page",
        "Adversaries may abuse Microsoft Outlook's Home Page feature to obtain persistence on a compromised system"
    ],
    "T1137.005": [
        "Outlook Rules",
        "Adversaries may abuse Microsoft Outlook rules to obtain persistence on a compromised system"
    ],
    "T1137.006": [
        "Add-ins",
        "Adversaries may abuse Microsoft Office add-ins to obtain persistence on a compromised system"
    ],
    "T1653": [
        "Power Settings",
        "Adversaries may impair a system's ability to hibernate, reboot, or shut down in order to extend access to infected machines"
    ],
    "T1542": [
        "Pre-OS Boot",
        "Adversaries may abuse Pre-OS Boot mechanisms as a way to establish persistence on a system"
    ],
    "T1542.001": [
        "System Firmware",
        "Adversaries may modify system firmware to persist on systems"
    ],
    "T1542.002": [
        "Component Firmware",
        "Adversaries may modify component firmware to persist on systems"
    ],
    "T1542.003": [
        "Bootkit",
        "Adversaries may use bootkits to persist on systems"
    ],
    "T1542.004": [
        "ROMMONkit",
        "Adversaries may abuse the ROM Monitor (ROMMON) by loading an unauthorized firmware with adversary code to provide persistent access and manipulate device behavior that is difficult to detect"
    ],
    "T1542.005": [
        "TFTP Boot",
        "Adversaries may abuse netbooting to load an unauthorized network device operating system from a Trivial File Transfer Protocol (TFTP) server"
    ],
    "T1053": [
        "Scheduled Task/Job",
        "Adversaries may abuse task scheduling functionality to facilitate initial or recurring execution of malicious code"
    ],
    "T1053.002": [
        "At",
        "Adversaries may abuse the at utility to perform task scheduling for initial or recurring execution of malicious code"
    ],
    "T1053.003": [
        "Cron",
        "Adversaries may abuse the cron utility to perform task scheduling for initial or recurring execution of malicious code"
    ],
    "T1053.005": [
        "Scheduled Task",
        "Adversaries may abuse the Windows Task Scheduler to perform task scheduling for initial or recurring execution of malicious code"
    ],
    "T1053.006": [
        "Systemd Timers",
        "Adversaries may abuse systemd timers to perform task scheduling for initial or recurring execution of malicious code"
    ],
    "T1053.007": [
        "Container Orchestration Job",
        "Adversaries may abuse task scheduling functionality provided by container orchestration tools such as Kubernetes to schedule deployment of containers configured to execute malicious code"
    ],
    "T1505": [
        "Server Software Component",
        "Adversaries may abuse legitimate extensible development features of servers to establish persistent access to systems"
    ]
}